# iso/airootfs/etc/skel/.config/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=Calamares Installer
Exec=calamares
X-GNOME-Autostart-enabled=true
